# cronjobs/crontab
*/5 * * * * php -f /var/www/BudgetSystem/cronjobs/register_accounts.php
*/10 * * * * php -f /var/www/BudgetSystem/cronjobs/process_transactions.php