# cronjobs/crontab
*/5 * * * * php -f /var/www/BudgetSystem/cronjobs/register_accounts.php